# systemd/midi-daemon-system.service
[Unit]
Description=MIDI Lua Routing Daemon (system-wide)
Documentation=https://github.com/yourname/midi-daemon
After=sound.target pipewire.service pipewire-pulse.service

[Service]
Type=simple
# Adjust the binary path for your distro:
#   - Arch/Debian/Ubuntu: /usr/local/bin/midi-daemon
#   - NixOS: use the Nix store path, e.g.
#       ${pkgs.midi-daemon}/bin/midi-daemon
#     or add midi-daemon to environment.systemPackages and use:
#       /run/current-system/sw/bin/midi-daemon
ExecStart=/usr/local/bin/midi-daemon \
    --config /etc/midi-daemon/config.toml \
    --routes /etc/midi-daemon/routes.d
ExecReload=/bin/kill -USR1 $MAINPID
Restart=on-failure
RestartSec=2
# NixOS: declare this user via users.users.midi-daemon in your configuration.nix,
# or replace with DynamicUser=yes (removes need for a pre-created user).
User=midi-daemon
Group=audio
KillSignal=SIGTERM
TimeoutStopSec=30
RuntimeDirectory=midi-daemon
CacheDirectory=midi-daemon

StandardOutput=journal
StandardError=journal
Environment=RUST_LOG=midi_daemon=info

[Install]
WantedBy=multi-user.target
