# systemd/midi-daemon-system.service
[Unit]
Description=MIDI Lua Routing Daemon (system-wide)
Documentation=https://github.com/yourname/midi-daemon
After=sound.target pipewire.service pipewire-pulse.service

[Service]
Type=simple
ExecStart=/usr/local/bin/midi-daemon
ExecReload=/bin/kill -USR1 $MAINPID
Restart=on-failure
RestartSec=2
User=midi-daemon
Group=midi-daemon
KillSignal=SIGTERM
TimeoutStopSec=30
RuntimeDirectory=midi-daemon
CacheDirectory=midi-daemon

StandardOutput=journal
StandardError=journal
Environment=RUST_LOG=midi_daemon=info

[Install]
WantedBy=multi-user.target
